# templates/projectopen.sh.erb
# MANAGED BY PUPPET

#!/bin/sh /etc/rc.common

/usr/local/aolserver451_rc2/bin/nsd -it <%= @serverroot %>/<%= @dbname %>/etc/config.tcl -u <%= @user %> -g <%= @group %> &

exit 0
